# config/displayboard-wifi.service
[Unit]
Description=DisplayBoard WiFi AP Fallback
After=network.target
Before=pm2-pi.service

[Service]
Type=forking
ExecStart=/home/pi/dakboard-local/scripts/ap-fallback.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
